# templates/default/main-service-1.conf.erb
start on starting <%= @name %>-<%= @process_name %>
stop on stopping <%= @name %>-<%= @process_name %>
respawn

exec su - <%= @owner %> -c 'cd <%= @path %>; <%= @command %> >> <%= @log_dir %>/<%= @name %>/<%= @process_name %>-1.log 2>&1'
